# systemd/admin_panel.service
[Unit]
Description=Admin Bots Panel
After=network.target postgresql.service

[Service]
Type=simple
User=adminbots
WorkingDirectory=/opt/admin-bots-platform
EnvironmentFile=/opt/admin-bots-platform/.env
Environment="PYTHONPATH=/opt/admin-bots-platform"
# Single-worker for faster startup (scales fine for small deployments)
ExecStart=/opt/admin-bots-platform/venv/bin/uvicorn admin_panel.app:app --host 127.0.0.1 --port 8000 --workers 1 --timeout-keep-alive 120
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=65535
LimitNPROC=4096
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
